# systemctl_config/hhinfo_main.service
[Unit]
Description=HHinfo Main Service


[Service]
Type=simple
WorkingDirectory=/home/ubuntu/hhinfo_PI
ExecStart=/home/ubuntu/hhinfo_PI/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target